# services/scxtop-mcp.service
[Unit]
Description=scxtop MCP daemon (sched_ext observability)

[Service]
ExecStart=/usr/bin/scxtop mcp --daemon
Restart=on-failure

# Hardening
PrivateNetwork=true
PrivateTmp=true
PrivateUsers=true
ProtectHome=true
ProtectSystem=strict
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_PERFMON
NoNewPrivileges=true
UMask=0077

[Install]
WantedBy=default.target
